# rules/legit-rules/proc_creation_win_schtasks_env_folder.yml
title: Suspicious Schtasks From Env Var Folder
id: 81325ce1-be01-4250-944f-b4789644556f
related:
- id: 43f487f0-755f-4c2a-bce7-d6d2eec2fcf8
  type: derived
status: experimental
description: Detects Schtask creations that point to a suspicious folder or an environment
  variable often used by malware
references:
- https://www.welivesecurity.com/2022/01/18/donot-go-do-not-respawn/
- https://www.joesandbox.com/analysis/514608/0/html#324415FF7D8324231381BAD48A052F85DF04
author: Florian Roth (Nextron Systems)
date: 2022/02/21
modified: 2023/11/30
tags:
- attack.execution
- attack.t1053.005
logsource:
  product: windows
  category: process_creation
detection:
- selection1_create:
  - and:
    - Image|endswith: \schtasks.exe
    - CommandLine|contains: ' /create '
- selection1_all_folders:
  - and:
    - or:
      - CommandLine|contains: :\Perflogs
      - CommandLine|contains: :\Windows\Temp
      - CommandLine|contains: \AppData\Local\
      - CommandLine|contains: \AppData\Roaming\
      - CommandLine|contains: \Users\Public
      - CommandLine|contains: '%AppData%'
      - CommandLine|contains: '%Public%'
- selection2_parent:
  - and:
    - ParentCommandLine|endswith: \svchost.exe -k netsvcs -p -s Schedule
- selection2_some_folders:
  - and:
    - or:
      - CommandLine|contains: :\Perflogs
      - CommandLine|contains: :\Windows\Temp
      - CommandLine|contains: \Users\Public
      - CommandLine|contains: '%Public%'
- filter_mixed:
  - and:
    - or:
      - or:
        - CommandLine|contains: update_task.xml
        - CommandLine|contains: /Create /TN TVInstallRestore /TR
      - ParentCommandLine|contains: unattended.ini
- filter_avira_install:
  - and:
    - CommandLine|contains: /Create /Xml "C:\Users\
    - CommandLine|contains: \AppData\Local\Temp\.CR.
    - CommandLine|contains: Avira_Security_Installation.xml
- filter_avira_other:
  - and:
    - CommandLine|contains: /Create /F /TN
    - CommandLine|contains: '/Xml '
    - CommandLine|contains: \AppData\Local\Temp\is-
    - CommandLine|contains: Avira_
    - or:
      - CommandLine|contains: .tmp\UpdateFallbackTask.xml
      - CommandLine|contains: .tmp\WatchdogServiceControlManagerTimeout.xml
      - CommandLine|contains: .tmp\SystrayAutostart.xml
      - CommandLine|contains: .tmp\MaintenanceTask.xml
- filter_klite_codec:
  - and:
    - CommandLine|contains: \AppData\Local\Temp\
    - CommandLine|contains: '/Create /TN "klcp_update" /XML '
    - CommandLine|contains: \klcp_update_task.xml
- condition: ( all of selection1* or all of selection2* ) and not 1 of filter*
falsepositives:
- Benign scheduled tasks creations or executions that happen often during software
  installations
- Software that uses the AppData folder and scheduled tasks to update the software
  in the AppData folders
level: medium
