# rules/legit-rules/proc_creation_win_schtasks_reg_loader_encoded.yml
title: Scheduled Task Executing Encoded Payload from Registry
id: c4eeeeae-89f4-43a7-8b48-8d1bdfa66c78
status: test
description: Detects the creation of a schtask that potentially executes a base64
  encoded payload stored in the Windows Registry using PowerShell.
references:
- https://thedfirreport.com/2022/02/21/qbot-and-zerologon-lead-to-full-domain-compromise/
author: pH-T (Nextron Systems), @Kostastsale, @TheDFIRReport, X__Junior (Nextron Systems),
  Nasreddine Bencherchali (Nextron Systems)
date: 2022/02/12
modified: 2023/02/04
tags:
- attack.execution
- attack.persistence
- attack.t1053.005
- attack.t1059.001
logsource:
  product: windows
  category: process_creation
detection:
- selection_img:
  - and:
    - or:
      - Image|endswith: \schtasks.exe
      - OriginalFileName|==: schtasks.exe
- selection_cli_create:
  - and:
    - CommandLine|contains: /Create
- selection_cli_encoding:
  - and:
    - or:
      - CommandLine|contains: FromBase64String
      - CommandLine|contains: encodedcommand
- selection_cli_get:
  - and:
    - or:
      - CommandLine|contains: Get-ItemProperty
      - CommandLine|contains: ' gp '
- selection_cli_hive:
  - and:
    - or:
      - CommandLine|contains: 'HKCU:'
      - CommandLine|contains: 'HKLM:'
      - CommandLine|contains: 'registry::'
      - CommandLine|contains: HKEY_
- condition: all of selection_*
falsepositives:
- Unlikely
level: high
