# rules/legit-rules/proc_creation_win_schtasks_susp_pattern.yml
title: Suspicious Command Patterns In Scheduled Task Creation
id: f2c64357-b1d2-41b7-849f-34d2682c0fad
status: experimental
description: Detects scheduled task creation using "schtasks" that contain potentially
  suspicious or uncommon commands
references:
- https://app.any.run/tasks/512c1352-6380-4436-b27d-bb62f0c020d6/
- https://twitter.com/RedDrip7/status/1506480588827467785
- https://www.ncsc.gov.uk/static-assets/documents/malware-analysis-reports/devil-bait/NCSC-MAR-Devil-Bait.pdf
author: Florian Roth (Nextron Systems)
date: 2022/02/23
modified: 2024/03/19
tags:
- attack.execution
- attack.t1053.005
logsource:
  product: windows
  category: process_creation
detection:
- and:
  - Image|endswith: \schtasks.exe
  - CommandLine|contains: '/Create '
  - or:
    - CommandLine|contains: '/sc minute '
    - CommandLine|contains: '/ru system '
  - or:
    - CommandLine|contains: cmd /c
    - CommandLine|contains: cmd /k
    - CommandLine|contains: cmd /r
    - CommandLine|contains: 'cmd.exe /c '
    - CommandLine|contains: 'cmd.exe /k '
    - CommandLine|contains: 'cmd.exe /r '
  - or:
    - CommandLine|contains: ' -decode '
    - CommandLine|contains: ' -enc '
    - CommandLine|contains: ' -w hidden '
    - CommandLine|contains: ' bypass '
    - CommandLine|contains: ' IEX'
    - CommandLine|contains: .DownloadData
    - CommandLine|contains: .DownloadFile
    - CommandLine|contains: .DownloadString
    - CommandLine|contains: '/c start /min '
    - CommandLine|contains: FromBase64String
    - CommandLine|contains: mshta http
    - CommandLine|contains: mshta.exe http
  - or:
    - CommandLine|contains: :\ProgramData\
    - CommandLine|contains: :\Temp\
    - CommandLine|contains: :\Tmp\
    - CommandLine|contains: :\Users\Public\
    - CommandLine|contains: :\Windows\Temp\
    - CommandLine|contains: \AppData\
    - CommandLine|contains: '%AppData%'
    - CommandLine|contains: '%Temp%'
    - CommandLine|contains: '%tmp%'
  - or:
    - CommandLine|contains: cscript
    - CommandLine|contains: curl
    - CommandLine|contains: wscript
falsepositives:
- Software installers that run from temporary folders and also install scheduled tasks
  are expected to generate some false positives
level: high
