# rules/active-rules/metadata/363eccc0-279a-4ccf-a3ab-24c2e63b11fb.yml
author: frack113
date: 2021/12/28
description: Adversaries may abuse the Windows Task Scheduler to perform task scheduling
  for initial or recurring execution of malicious code
falsepositives:
- Unknown
level: medium
references:
- https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1053.005/T1053.005.md#atomic-test-4---powershell-cmdlet-scheduled-task
- https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1053.005/T1053.005.md#atomic-test-6---wmi-invoke-cimmethod-scheduled-task
status: test
tags:
- attack.persistence
- attack.t1053.005
title: Powershell Create Scheduled Task
