# rules/active-rules/metadata/39019a4e-317f-4ce3-ae63-309a8c6b53c5.yml
author: Florian Roth (Nextron Systems)
date: 2021/03/11
description: Detects the creation of scheduled tasks that involves a temporary folder
  and runs only once
falsepositives:
- Administrative activity
- Software installation
fields:
- CommandLine
- ParentCommandLine
level: high
modified: 2022/10/09
references:
- https://discuss.elastic.co/t/detection-and-response-for-hafnium-activity/266289/3
status: test
tags:
- attack.execution
- attack.persistence
- attack.t1053.005
title: Suspicious Scheduled Task Creation Involving Temp Folder
