# rules/active-rules/metadata/424273ea-7cf8-43a6-b712-375f925e481f.yml
author: Nasreddine Bencherchali (Nextron Systems)
date: 2022/12/05
description: Detects the execution of Scheduled Tasks where the Program being run
  is located in a suspicious location or it's an unusale program to be run from a
  Scheduled Task
falsepositives:
- Unknown
level: medium
modified: 2023/02/07
references:
- Internal Research
status: test
tags:
- attack.persistence
- attack.t1053.005
title: Scheduled Task Executed From A Suspicious Location
