# rules/active-rules/metadata/4720b7df-40c3-48fd-bbdf-fd4b3c464f0d.yml
author: Syed Hasan (@syedhasan009)
date: 2021/06/18
description: Monitor the creation of a new key under 'TaskCache' when a new scheduled
  task is registered by a process that is not svchost.exe, which is suspicious
falsepositives:
- Unknown
level: high
modified: 2023/08/17
references:
- https://thedfirreport.com/2021/03/29/sodinokibi-aka-revil-ransomware/
- https://labs.f-secure.com/blog/scheduled-task-tampering/
status: experimental
tags:
- attack.persistence
- attack.t1053
- attack.t1053.005
title: Scheduled TaskCache Change by Uncommon Program
