# rules/active-rules/metadata/614cf376-6651-47c4-9dcc-6b9527f749f4.yml
author: Nasreddine Bencherchali (Nextron Systems)
date: 2022/12/05
description: Detects update to a scheduled task event that contain suspicious keywords.
falsepositives:
- Unknown
level: high
references:
- https://learn.microsoft.com/en-us/windows/security/threat-protection/auditing/event-4698
related:
- id: 1c0e41cd-21bb-4433-9acc-4a2cd6367b9b
  type: similar
status: test
tags:
- attack.execution
- attack.privilege_escalation
- attack.persistence
- attack.t1053.005
title: Suspicious Scheduled Task Update
