# rules/active-rules/metadata/80e1f67a-4596-4351-98f5-a9c3efabac95.yml
author: Florian Roth (Nextron Systems)
date: 2021/11/16
description: Detects the creation of tasks from processes executed from suspicious
  locations
falsepositives:
- Unknown
level: high
modified: 2022/01/12
references:
- Internal Research
status: test
tags:
- attack.persistence
- attack.execution
- attack.t1053
title: Suspicious Scheduled Task Write to System32 Tasks
