# rules/active-rules/metadata/8a8379b8-780b-4dbf-b1e9-31c8d112fefb.yml
author: Florian Roth (Nextron Systems)
date: 2022/04/15
description: Detects scheduled task creations that have suspicious action command
  and folder combinations
falsepositives:
- Unknown
level: high
modified: 2022/11/18
references:
- https://symantec-enterprise-blogs.security.com/blogs/threat-intelligence/lazarus-dream-job-chemical
status: test
tags:
- attack.execution
- attack.t1053.005
title: Schtasks From Suspicious Folders
