# rules/active-rules/metadata/92626ddd-662c-49e3-ac59-f6535f12d189.yml
author: Florian Roth (Nextron Systems)
date: 2019/01/16
description: Detects the creation of scheduled tasks by user accounts via the "schtasks"
  utility.
falsepositives:
- Administrative activity
- Software installation
level: low
modified: 2024/01/18
references:
- https://learn.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks-create
status: test
tags:
- attack.execution
- attack.persistence
- attack.privilege_escalation
- attack.t1053.005
- attack.s0111
- car.2013-08-001
- stp.1u
title: Scheduled Task Creation Via Schtasks.EXE
