# rules/active-rules/metadata/970823b7-273b-460a-8afc-3a6811998529.yml
author: pH-T (Nextron Systems)
date: 2022/07/15
description: Detects scheduled task creation events that include suspicious actions,
  and is run once at 00:00
falsepositives:
- Software installation
filename: proc_creation_win_schtasks_one_time_only_midnight_task.yml
level: high
logsource:
  category: process_creation
  product: windows
modified: 2023/02/03
references:
- https://www.trendmicro.com/vinfo/us/security/news/ransomware-spotlight/ransomware-spotlight-blackbyte
related: []
status: test
tags:
- attack.execution
- attack.persistence
- attack.privilege_escalation
- attack.t1053.005
