# rules/active-rules/metadata/a8f29a7b-b137-4446-80a0-b804272f3da2.yml
author: Samir Bousseaden
date: 2019/04/03
description: Detect lateral movement using GPO scheduled task, usually used to deploy
  ransomware at scale
falsepositives:
- If the source IP is not localhost then it's super suspicious, better to monitor
  both local and remote changes to GPO scheduledtasks
level: high
modified: 2021/11/27
references:
- https://twitter.com/menasec1/status/1106899890377052160
- https://www.secureworks.com/blog/ransomware-as-a-distraction
status: test
tags:
- attack.persistence
- attack.lateral_movement
- attack.t1053.005
title: Persistence and Execution at Scale via GPO Scheduled Task
