# rules/active-rules/metadata/c4eeeeae-89f4-43a7-8b48-8d1bdfa66c78.yml
author: pH-T (Nextron Systems), @Kostastsale, @TheDFIRReport, X__Junior (Nextron Systems),
  Nasreddine Bencherchali (Nextron Systems)
date: 2022/02/12
description: Detects the creation of a schtask that potentially executes a base64
  encoded payload stored in the Windows Registry using PowerShell.
falsepositives:
- Unlikely
filename: proc_creation_win_schtasks_reg_loader_encoded.yml
level: high
logsource:
  category: process_creation
  product: windows
modified: 2023/02/04
references:
- https://thedfirreport.com/2022/02/21/qbot-and-zerologon-lead-to-full-domain-compromise/
related: []
status: test
tags:
- attack.execution
- attack.persistence
- attack.t1053.005
- attack.t1059.001
