# rules/active-rules/metadata/c5c00f49-b3f9-45a6-997e-cfdecc6e1967.yml
author: pH-T (Nextron Systems), Nasreddine Bencherchali (Nextron Systems)
date: 2022/03/15
description: Detects the creation of a schtask that executes a file from C:\Users\<USER>\AppData\Local
falsepositives:
- Unknown
filename: proc_creation_win_schtasks_appdata_local_system.yml
level: high
modified: 2022/07/28
references:
- https://thedfirreport.com/2022/02/21/qbot-and-zerologon-lead-to-full-domain-compromise/
related: []
status: test
tags:
- attack.execution
- attack.persistence
- attack.t1053.005
- attack.t1059.001
