# rules/active-rules/metadata/dd2a821e-3b07-4d3b-a9ac-929fe4c6ca0c.yml
author: Swachchhanda Shrawan Poudel, Elastic (idea)
date: 2023/04/20
description: Detects the creation of a scheduled task using the "-XML" flag with a
  file without the '.xml' extension. This behavior could be indicative of potential
  defense evasion attempt during persistence
falsepositives:
- Unknown
level: medium
references:
- https://docs.microsoft.com/en-us/windows/win32/taskschd/daily-trigger-example--xml-
- https://github.com/elastic/protections-artifacts/blob/084067123d3328a823b1c3fdde305b694275c794/behavior/rules/persistence_suspicious_scheduled_task_creation_via_masqueraded_xml_file.toml
status: test
tags:
- attack.defense_evasion
- attack.persistence
- attack.t1036.005
- attack.t1053.005
title: Suspicious Scheduled Task Creation via Masqueraded XML File
