# rules/active-rules/metadata/ff2fff64-4cd6-4a2b-ba7d-e28a30bbe66b.yml
author: Nasreddine Bencherchali (Nextron Systems)
date: 2022/10/31
description: Detects creation of a scheduled task with a GUID like name
falsepositives:
- Legitimate software naming their tasks as GUIDs
filename: proc_creation_win_schtasks_guid_task_name.yml
level: medium
modified: ''
references:
- https://thedfirreport.com/2022/10/31/follina-exploit-leads-to-domain-compromise/
- https://thedfirreport.com/2022/02/21/qbot-and-zerologon-lead-to-full-domain-compromise/
related: []
status: test
tags:
- attack.execution
- attack.t1053.005
